# packages/vdr/vdr-depends/_vtuner-ng/_system.d/vtuner-ng.service
[Unit]
Description=vtuner-ng
Requires=network-online.target graphical.target

[Service]
Type=forking
ExecStop=/usr/bin/bash /storage/.config/stop_vtuner.sh
ExecStart=/usr/bin/bash /storage/.config/start_vtuner.sh

[Install]
WantedBy=multi-user.target

