# packages/vdr/_config/system.d/VDRCoreELEC.service
[Unit]
Description=Preparing VDRCoreELEC
Before=kodi.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/VDRCoreELEC.sh

[Install]
WantedBy=multi-user.target