# packages/vdr/_vdr/_system.d/switch_kodi_vdr.service
[Service]  
Type=simple  
ExecStart=XXBINDIRXX/switch_kodi_vdr.sh

[Install]
WantedBy=multi-user.target
