# packages/vdr/_vdr/_system.d/switch_kodi_vdr.path
[Unit]
Wants=switch_kodi_vdr.service

[Path]
PathChanged=/storage/.cache/switch_kodi_vdr

[Install]
WantedBy=multi-user.target
