# packages/vdr/_vdr/_system.d/vdropt.target
[Unit]
Description=VDR Interface
Requires=multi-user.target graphical.target network-online.target
After=network-online.target graphical.target
Conflicts=rescue.target
AllowIsolate=yes

[Install]
Alias=vdropt.target
