# packages/vdr/_vdr/_system.d/vdropt.service
[Unit]
Description=VDR Client
Requires=network-online.target graphical.target
After=network-online.target graphical.target

[Service]
ExecStart=/usr/local/bin/start_vdr.sh
Restart=always
RestartSec=2
StartLimitInterval=0
TimeoutStartSec=30
TimeoutStopSec=30
LimitNOFILE=16384

[Install]
WantedBy=vdropt.target

