# packages/vdr/vdr-depends/_remotetranscode/system.d/remotetranscode.service
[Unit]
Description=remotetranscode
Requires=network-online.target
After=network-online.target
StartLimitIntervalSec=400
StartLimitBurst=5

[Service]
WorkingDirectory=/storage/remotetranscode
ExecStart=/storage/remotetranscode/remotrans -c /storage/.config/vdropt/sockets.ini
Restart=always
RestartSec=30
StartLimitInterval=0
TimeoutStartSec=30
TimeoutStopSec=30
LimitNOFILE=16384

[Install]
WantedBy=vdropt.target
