# packages/vdr/_config/_system.d/switch_kodi_vdr.path
[Unit]
Description=Service welcher auf Dateiänderung reagiert

[Path]
PathChanged=/storage/.cache/switch_kodi_vdr
Unit=switch_kodi_vdr.service

[Install]
WantedBy=multi-user.target