# packages/vdr/vdr-depends/_spice_vdagent/system.d/spice-vdagentd.socket
[Unit]
Description=Activation socket for spice guest agent daemon

[Socket]
ListenStream=/run/spice-vdagentd/spice-vdagent-sock
