# packages/vdr/vdr-depends/_spice_vdagent/system.d/spice-vdagent.service
[Unit]
Description=Spice guest session agent

[Service]
ExecStart=/usr/bin/spice-vdagent -x
